# deploy/proxmox-tracker-manager.service
[Unit]
Description=Proxmox VM Time Tracker Manager
After=network.target mysql.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/proxmox-tracker

# Use the Python from venv - update this path if different
ExecStart=/opt/proxmox-tracker/venv/bin/python -m uvicorn manager.main:app --host 0.0.0.0 --port 8000

Restart=always
RestartSec=10

# Environment
Environment=PYTHONUNBUFFERED=1

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=proxmox-tracker-manager

[Install]
WantedBy=multi-user.target
